# app/crontab
*/10 * * * * python /speed_monitor/main.py measure >> /logs/cron.log